# stdlib/supervisord.conf
[supervisord] 
nodaemon=true

[program:sshd] 
# command=/usr/sbin/sshd -D
command=service ssh start

[program:apache] 
command=/usr/sbin/apache2ctl -D FOREGROUND